# systemd/user/prima-qwen397b-babysitter.timer
[Unit]
Description=Run the PRIMA Qwen3.5 397B vertical-line babysitter every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
