# systemd/user/prima-qwen397b-babysitter.timer
[Unit]
Description=Run the PRIMA Qwen3.5 397B vertical-line babysitter hourly

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
