# systemd/user/prima-qwen397b-babysitter.timer
[Unit]
Description=Run the PRIMA Qwen3.5 397B vertical-line babysitter every 30 minutes

[Timer]
OnCalendar=*:0/30
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
